# docker/apache/poem-apache.ini
[supervisord]
nodaemon=true

[program:httpd]
command=httpd
autorestart=true
