# testenv/apache/poem-apache.ini
[supervisord]
nodaemon=true

[program:httpd]
command=scl enable httpd24 'httpd'
autorestart=true
